# cluster/pgsql/playbooks/roles/base_config/templates/patroni.service.j2
[Unit]
Description=Patroni PostgreSQL high availability
After=network-online.target etcd.service
Wants=network-online.target etcd.service

[Service]
Type=simple
User=postgres
Group=postgres
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# /run/postgresql lives on tmpfs and vanishes on every reboot (including a
# stop/resize/start vertical scale). Without it Postgres dies at startup
# creating its socket lock, so Patroni crash-loops and start/add-member hang
# waiting for the port. Let systemd own the directory: recreated on every
# Patroni start, preserved across restarts (KillMode=process leaves postgres
# running through a Patroni restart, so the socket dir must survive it too).
RuntimeDirectory=postgresql
RuntimeDirectoryMode=2775
RuntimeDirectoryPreserve=yes
ExecStart={{ erawan_patroni_bin_path }} {{ patroni_config_path }}
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=process
TimeoutSec=30
Restart=on-failure
RestartSec=5
LimitNOFILE=262144

[Install]
WantedBy=multi-user.target
