# cluster/mysql/playbooks/roles/auto_rejoin/templates/gr_watchdog.service.j2
[Unit]
Description=Erawan GR watchdog — auto-heal stuck Group Replication applier
After=mysql.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/erawan-gr-watchdog
